# scripts/cloneclean.service.j2
#!/bin/bash
#
################################################################################
[Unit]
Description=cloneclean service
After=multi-user.target network-online.target

[Service]
Type=simple
User=1000
Group=1000
ExecStart=/bin/bash /opt/pgclone/scripts/cloneclean.sh
KillMode=process

[Install]
WantedBy=multi-user.target
